# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules-threat-hunting/windows/registry/registry_event/registry_event_scheduled_task_creation.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules-threat-hunting/windows/registry/registry_event/registry_event_scheduled_task_creation.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Center for Threat Informed Defense (CTID) Summiting the Pyramid Team",
                "date": "2023-09-27",
                "title": "Scheduled Task Created - Registry",
                "id": "93ff0ceb-e0ef-4586-8cd8-a6c277d738e3"
            },
            "rule_text": "title: Scheduled Task Created - Registry\nid: 93ff0ceb-e0ef-4586-8cd8-a6c277d738e3\nstatus: test\ndescription: Detects the creation of a scheduled task via Registry keys.\nreferences:\n- https://center-for-threat-informed-defense.github.io/summiting-the-pyramid/analytics/task_scheduling/\n- https://posts.specterops.io/abstracting-scheduled-tasks-3b6451f6a1c5\nauthor: Center for Threat Informed Defense (CTID) Summiting the Pyramid Team\ndate: 2023-09-27\ntags:\n- attack.execution\n- attack.persistence\n- attack.privilege-escalation\n- attack.s0111\n- attack.t1053.005\n- car.2013-08-001\n- detection.threat-hunting\nlogsource:\n  product: windows\n  category: registry_event\ndetection:\n  selection:\n    TargetObject|contains:\n    - \\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tasks\\\n    - \\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree\\\n  condition: selection\nfalsepositives:\n- Likely as this is a normal behaviour on Windows\nlevel: low\n",
            "vulnerabilities": []
        }
    ]
}