# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/macos/process_creation/proc_creation_macos_schedule_task_job_cron.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/macos/process_creation/proc_creation_macos_schedule_task_job_cron.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Alejandro Ortuno, oscd.community",
                "date": "2020-10-06",
                "title": "Scheduled Cron Task/Job - MacOs",
                "id": "7c3b43d8-d794-47d2-800a-d277715aa460"
            },
            "rule_text": "title: Scheduled Cron Task/Job - MacOs\nid: 7c3b43d8-d794-47d2-800a-d277715aa460\nstatus: test\ndescription: Detects abuse of the cron utility to perform task scheduling for initial\n  or recurring execution of malicious code. Detection will focus on crontab jobs uploaded\n  from the tmp folder.\nreferences:\n- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.003/T1053.003.md\nauthor: Alejandro Ortuno, oscd.community\ndate: 2020-10-06\nmodified: 2022-11-27\ntags:\n- attack.execution\n- attack.persistence\n- attack.privilege-escalation\n- attack.t1053.003\nlogsource:\n  category: process_creation\n  product: macos\ndetection:\n  selection:\n    Image|endswith: /crontab\n    CommandLine|contains: /tmp/\n  condition: selection\nfalsepositives:\n- Legitimate administration activities\nlevel: medium\n",
            "vulnerabilities": []
        }
    ]
}