# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules-emerging-threats/2018/TA/Slingshot/proc_creation_win_apt_slingshot.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules-emerging-threats/2018/TA/Slingshot/proc_creation_win_apt_slingshot.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Florian Roth (Nextron Systems), Bartlomiej Czyz (@bczyz1)",
                "date": "2019-03-04",
                "title": "Defrag Deactivation",
                "id": "958d81aa-8566-4cea-a565-59ccd4df27b0"
            },
            "rule_text": "title: Defrag Deactivation\nid: 958d81aa-8566-4cea-a565-59ccd4df27b0\nstatus: test\ndescription: Detects the deactivation and disabling of the Scheduled defragmentation\n  task as seen by Slingshot APT group\nreferences:\n- https://securelist.com/apt-slingshot/84312/\nauthor: Florian Roth (Nextron Systems), Bartlomiej Czyz (@bczyz1)\ndate: 2019-03-04\nmodified: 2022-10-09\ntags:\n- attack.privilege-escalation\n- attack.execution\n- attack.persistence\n- attack.t1053.005\n- attack.s0111\n- detection.emerging-threats\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: \\schtasks.exe\n    CommandLine|contains:\n    - /delete\n    - /change\n    CommandLine|contains|all:\n    - /TN\n    - \\Microsoft\\Windows\\Defrag\\ScheduledDefrag\n  condition: selection\nfalsepositives:\n- Unknown\nlevel: medium\n",
            "vulnerabilities": []
        }
    ]
}