# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/registry/registry_set/registry_set_timeproviders_dllname.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/registry/registry_set/registry_set_timeproviders_dllname.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "frack113",
                "date": "2022-06-19",
                "title": "New TimeProviders Registered With Uncommon DLL Name",
                "id": "e88a6ddc-74f7-463b-9b26-f69fc0d2ce85"
            },
            "rule_text": "title: New TimeProviders Registered With Uncommon DLL Name\nid: e88a6ddc-74f7-463b-9b26-f69fc0d2ce85\nstatus: test\ndescription: 'Detects processes setting a new DLL in DllName in under HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Services\\W32Time\\TimeProvider.\n\n  Adversaries may abuse time providers to execute DLLs when the system boots.\n\n  The Windows Time service (W32Time) enables time synchronization across and within\n  domains.\n\n  '\nreferences:\n- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.003/T1547.003.md\nauthor: frack113\ndate: 2022-06-19\nmodified: 2024-03-26\ntags:\n- attack.persistence\n- attack.privilege-escalation\n- attack.t1547.003\nlogsource:\n  category: registry_set\n  product: windows\ndetection:\n  selection:\n    TargetObject|contains: \\Services\\W32Time\\TimeProviders\n    TargetObject|endswith: \\DllName\n  filter_main_w32time:\n    Details:\n    - '%SystemRoot%\\System32\\vmictimeprovider.dll'\n    - '%systemroot%\\system32\\w32time.dll'\n    - C:\\Windows\\SYSTEM32\\w32time.DLL\n  condition: selection and not 1 of filter_main_*\nfalsepositives:\n- Unknown\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}