# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/builtin/security/win_security_susp_scheduled_task_delete_or_disable.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/builtin/security/win_security_susp_scheduled_task_delete_or_disable.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Nasreddine Bencherchali (Nextron Systems)",
                "date": "2022-12-05",
                "title": "Important Scheduled Task Deleted/Disabled",
                "id": "7595ba94-cf3b-4471-aa03-4f6baa9e5fad"
            },
            "rule_text": "title: Important Scheduled Task Deleted/Disabled\nid: 7595ba94-cf3b-4471-aa03-4f6baa9e5fad\nrelated:\n- id: dbc1f800-0fe0-4bc0-9c66-292c2abe3f78\n  type: similar\n- id: 9ac94dc8-9042-493c-ba45-3b5e7c86b980\n  type: similar\n- id: 9e3cb244-bdb8-4632-8c90-6079c8f4f16d\n  type: similar\nstatus: test\ndescription: Detects when adversaries stop services or processes by deleting or disabling\n  their respective scheduled tasks in order to conduct data destructive activities\nreferences:\n- https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4699\n- https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4701\nauthor: Nasreddine Bencherchali (Nextron Systems)\ndate: 2022-12-05\nmodified: 2023-03-13\ntags:\n- attack.execution\n- attack.privilege-escalation\n- attack.persistence\n- attack.t1053.005\nlogsource:\n  product: windows\n  service: security\n  definition: The Advanced Audit Policy setting Object Access > Audit Other Object\n    Access Events has to be configured to allow this detection. We also recommend\n    extracting the Command field from the embedded XML in the event data.\ndetection:\n  selection:\n    EventID:\n    - 4699\n    - 4701\n    TaskName|contains:\n    - \\Windows\\SystemRestore\\SR\n    - \\Windows\\Windows Defender\\\n    - \\Windows\\BitLocker\n    - \\Windows\\WindowsBackup\\\n    - \\Windows\\WindowsUpdate\\\n    - \\Windows\\UpdateOrchestrator\\Schedule\n    - \\Windows\\ExploitGuard\n  filter_sys_username:\n    EventID: 4699\n    SubjectUserName|endswith: $\n    TaskName|contains: \\Windows\\Windows Defender\\\n  condition: selection and not 1 of filter_*\nfalsepositives:\n- Unknown\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}