# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/builtin/security/win_security_susp_scheduled_task_creation.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/builtin/security/win_security_susp_scheduled_task_creation.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Nasreddine Bencherchali (Nextron Systems)",
                "date": "2022-12-05",
                "title": "Suspicious Scheduled Task Creation",
                "id": "3a734d25-df5c-4b99-8034-af1ddb5883a4"
            },
            "rule_text": "title: Suspicious Scheduled Task Creation\nid: 3a734d25-df5c-4b99-8034-af1ddb5883a4\nstatus: test\ndescription: Detects suspicious scheduled task creation events. Based on attributes\n  such as paths, commands line flags, etc.\nreferences:\n- https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4698\nauthor: Nasreddine Bencherchali (Nextron Systems)\ndate: 2022-12-05\nmodified: 2022-12-07\ntags:\n- attack.execution\n- attack.privilege-escalation\n- attack.persistence\n- attack.t1053.005\nlogsource:\n  product: windows\n  service: security\n  definition: The Advanced Audit Policy setting Object Access > Audit Other Object\n    Access Events has to be configured to allow this detection. We also recommend\n    extracting the Command field from the embedded XML in the event data.\ndetection:\n  selection_eid:\n    EventID: 4698\n  selection_paths:\n    TaskContent|contains:\n    - \\AppData\\Local\\Temp\\\n    - \\AppData\\Roaming\\\n    - \\Users\\Public\\\n    - \\WINDOWS\\Temp\\\n    - C:\\Temp\\\n    - \\Desktop\\\n    - \\Downloads\\\n    - \\Temporary Internet\n    - C:\\ProgramData\\\n    - C:\\Perflogs\\\n  selection_commands:\n    TaskContent|contains:\n    - regsvr32\n    - rundll32\n    - cmd.exe</Command>\n    - cmd</Command>\n    - '<Arguments>/c '\n    - '<Arguments>/k '\n    - '<Arguments>/r '\n    - powershell\n    - pwsh\n    - mshta\n    - wscript\n    - cscript\n    - certutil\n    - bitsadmin\n    - bash.exe\n    - 'bash '\n    - scrcons\n    - 'wmic '\n    - wmic.exe\n    - forfiles\n    - scriptrunner\n    - hh.exe\n  condition: all of selection_*\nfalsepositives:\n- Unknown\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}