# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/builtin/security/win_security_gpo_scheduledtasks.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/builtin/security/win_security_gpo_scheduledtasks.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Samir Bousseaden",
                "date": "2019-04-03",
                "title": "Persistence and Execution at Scale via GPO Scheduled Task",
                "id": "a8f29a7b-b137-4446-80a0-b804272f3da2"
            },
            "rule_text": "title: Persistence and Execution at Scale via GPO Scheduled Task\nid: a8f29a7b-b137-4446-80a0-b804272f3da2\nstatus: test\ndescription: Detect lateral movement using GPO scheduled task, usually used to deploy\n  ransomware at scale\nreferences:\n- https://twitter.com/menasec1/status/1106899890377052160\n- https://www.secureworks.com/blog/ransomware-as-a-distraction\n- https://www.elastic.co/guide/en/security/7.17/prebuilt-rule-0-16-1-scheduled-task-execution-at-scale-via-gpo.html\nauthor: Samir Bousseaden\ndate: 2019-04-03\nmodified: 2024-09-04\ntags:\n- attack.privilege-escalation\n- attack.execution\n- attack.persistence\n- attack.lateral-movement\n- attack.t1053.005\nlogsource:\n  product: windows\n  service: security\n  definition: The advanced audit policy setting \"Object Access > Audit Detailed File\n    Share\" must be configured for Success/Failure\ndetection:\n  selection_5136:\n    EventID: 5136\n    AttributeLDAPDisplayName:\n    - gPCMachineExtensionNames\n    - gPCUserExtensionNames\n    AttributeValue|contains:\n    - CAB54552-DEEA-4691-817E-ED4A4D1AFC72\n    - AADCED64-746C-4633-A97C-D61349046527\n  selection_5145:\n    EventID: 5145\n    ShareName|endswith: \\SYSVOL\n    RelativeTargetName|endswith: ScheduledTasks.xml\n    AccessList|contains:\n    - WriteData\n    - '%%4417'\n  condition: 1 of selection_*\nfalsepositives:\n- If the source IP is not localhost then it's super suspicious, better to monitor\n  both local and remote changes to GPO scheduled tasks.\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}