# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/registry/registry_event/registry_event_susp_lsass_dll_load.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/registry/registry_event/registry_event_susp_lsass_dll_load.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Florian Roth (Nextron Systems)",
                "date": "2019-10-16",
                "title": "DLL Load via LSASS",
                "id": "b3503044-60ce-4bf4-bbcb-e3db98788823"
            },
            "rule_text": "title: DLL Load via LSASS\nid: b3503044-60ce-4bf4-bbcb-e3db98788823\nstatus: test\ndescription: Detects a method to load DLL via LSASS process using an undocumented\n  Registry key\nreferences:\n- https://blog.xpnsec.com/exploring-mimikatz-part-1/\n- https://twitter.com/SBousseaden/status/1183745981189427200\nauthor: Florian Roth (Nextron Systems)\ndate: 2019-10-16\nmodified: 2022-04-21\ntags:\n- attack.privilege-escalation\n- attack.execution\n- attack.persistence\n- attack.t1547.008\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    TargetObject|contains:\n    - \\CurrentControlSet\\Services\\NTDS\\DirectoryServiceExtPt\n    - \\CurrentControlSet\\Services\\NTDS\\LsaDbExtPt\n  filter_domain_controller:\n    Image: C:\\Windows\\system32\\lsass.exe\n    Details:\n    - '%%systemroot%%\\system32\\ntdsa.dll'\n    - '%%systemroot%%\\system32\\lsadb.dll'\n  condition: selection and not 1 of filter_*\nfalsepositives:\n- Unknown\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}