# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules-emerging-threats/2020/TA/Leviathan/registry_event_apt_leviathan.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules-emerging-threats/2020/TA/Leviathan/registry_event_apt_leviathan.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Aidan Bracher",
                "date": "2020-07-07",
                "title": "Leviathan Registry Key Activity",
                "id": "70d43542-cd2d-483c-8f30-f16b436fd7db"
            },
            "rule_text": "title: Leviathan Registry Key Activity\nid: 70d43542-cd2d-483c-8f30-f16b436fd7db\nstatus: test\ndescription: Detects registry key used by Leviathan APT in Malaysian focused campaign\nreferences:\n- https://www.elastic.co/blog/advanced-techniques-used-in-malaysian-focused-apt-campaign\nauthor: Aidan Bracher\ndate: 2020-07-07\nmodified: 2023-09-19\ntags:\n- attack.privilege-escalation\n- attack.persistence\n- attack.t1547.001\n- detection.emerging-threats\nlogsource:\n  category: registry_event\n  product: windows\ndetection:\n  selection:\n    TargetObject|contains: \\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\ntkd\n  condition: selection\nlevel: critical\n",
            "vulnerabilities": []
        }
    ]
}