# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/linux/builtin/cron/lnx_cron_crontab_file_modification.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/linux/builtin/cron/lnx_cron_crontab_file_modification.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Pawel Mazur",
                "date": "2022-04-16",
                "title": "Modifying Crontab",
                "id": "af202fd3-7bff-4212-a25a-fb34606cfcbe"
            },
            "rule_text": "title: Modifying Crontab\nid: af202fd3-7bff-4212-a25a-fb34606cfcbe\nstatus: test\ndescription: Detects suspicious modification of crontab file.\nreferences:\n- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.003/T1053.003.md\nauthor: Pawel Mazur\ndate: 2022-04-16\ntags:\n- attack.privilege-escalation\n- attack.execution\n- attack.persistence\n- attack.t1053.003\nlogsource:\n  product: linux\n  service: cron\ndetection:\n  keywords:\n  - REPLACE\n  condition: keywords\nfalsepositives:\n- Legitimate modification of crontab\nlevel: medium\n",
            "vulnerabilities": []
        }
    ]
}