# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_schtasks_creation.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_schtasks_creation.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Florian Roth (Nextron Systems)",
                "date": "2019-01-16",
                "title": "Scheduled Task Creation Via Schtasks.EXE",
                "id": "92626ddd-662c-49e3-ac59-f6535f12d189"
            },
            "rule_text": "title: Scheduled Task Creation Via Schtasks.EXE\nid: 92626ddd-662c-49e3-ac59-f6535f12d189\nstatus: test\ndescription: Detects the creation of scheduled tasks by user accounts via the \"schtasks\"\n  utility.\nreferences:\n- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create\nauthor: Florian Roth (Nextron Systems)\ndate: 2019-01-16\nmodified: 2025-10-22\ntags:\n- attack.execution\n- attack.persistence\n- attack.privilege-escalation\n- attack.t1053.005\n- attack.s0111\n- car.2013-08-001\n- stp.1u\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: \\schtasks.exe\n    CommandLine|contains: ' /create '\n  filter_main_system_user:\n    User|contains:\n    - AUTHORI\n    - AUTORI\n  filter_optional_msoffice:\n    ParentImage:\n    - C:\\Program Files\\Microsoft Office\\root\\integration\\integrator.exe\n    - C:\\Program Files (x86)\\Microsoft Office\\root\\integration\\integrator.exe\n    Image:\n    - C:\\Windows\\System32\\schtasks.exe\n    - C:\\Windows\\SysWOW64\\schtasks.exe\n    CommandLine|contains: Microsoft\\Office\\Office Performance Monitor\n  condition: selection and not 1 of filter_main_* and not 1 of filter_optional_*\nfalsepositives:\n- Administrative activity\n- Software installation\nlevel: low\n",
            "vulnerabilities": []
        }
    ]
}