# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_renamed_schtasks_execution.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_renamed_schtasks_execution.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "experimental",
                "author": "Swachchhanda Shrawan Poudel (Nextron Systems)",
                "date": "2025-11-27",
                "title": "Renamed Schtasks Execution",
                "id": "f91e51c9-f344-4b32-969b-0b6f6b8537d4"
            },
            "rule_text": "title: Renamed Schtasks Execution\nid: f91e51c9-f344-4b32-969b-0b6f6b8537d4\nstatus: experimental\ndescription: 'Detects the execution of renamed schtasks.exe binary, which is a legitimate\n  Windows utility used for scheduling tasks.\n\n  One of the very common persistence techniques is schedule malicious tasks using\n  schtasks.exe.\n\n  Since, it is heavily abused, it is also heavily monitored by security products.\n  To evade detection, threat actors may rename the schtasks.exe binary to schedule\n  their malicious tasks.\n\n  '\nreferences:\n- https://x.com/JangPr0/status/1932034543026065833\n- https://ss64.com/nt/schtasks.html\nauthor: Swachchhanda Shrawan Poudel (Nextron Systems)\ndate: 2025-11-27\ntags:\n- attack.defense-evasion\n- attack.execution\n- attack.persistence\n- attack.privilege-escalation\n- attack.t1036.003\n- attack.t1053.005\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_cmd_operation:\n    CommandLine|contains|windash:\n    - ' /create '\n    - ' /delete '\n    - ' /query '\n    - ' /change '\n    - ' /run '\n    - ' /end '\n  selection_cmd_flags:\n    CommandLine|contains|windash:\n    - ' /tn '\n    - ' /tr '\n    - ' /sc '\n    - ' /st '\n    - ' /ru '\n    - ' /fo '\n  selection_pe:\n    OriginalFileName: schtasks.exe\n  filter_main_cmd:\n    CommandLine|contains: schtasks\n  filter_main_img:\n    Image|endswith: \\schtasks.exe\n  condition: (all of selection_cmd_* and not filter_main_cmd) or (selection_pe and\n    not filter_main_img)\nfalsepositives:\n- Unlikely\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}