# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_schtasks_env_folder.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_schtasks_env_folder.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Florian Roth (Nextron Systems)",
                "date": "2022-02-21",
                "title": "Schedule Task Creation From Env Variable Or Potentially Suspicious Path Via Schtasks.EXE",
                "id": "81325ce1-be01-4250-944f-b4789644556f"
            },
            "rule_text": "title: Schedule Task Creation From Env Variable Or Potentially Suspicious Path Via\n  Schtasks.EXE\nid: 81325ce1-be01-4250-944f-b4789644556f\nrelated:\n- id: 43f487f0-755f-4c2a-bce7-d6d2eec2fcf8\n  type: derived\nstatus: test\ndescription: Detects Schtask creations that point to a suspicious folder or an environment\n  variable often used by malware\nreferences:\n- https://www.welivesecurity.com/2022/01/18/donot-go-do-not-respawn/\n- https://www.joesandbox.com/analysis/514608/0/html#324415FF7D8324231381BAD48A052F85DF04\n- https://blog.talosintelligence.com/gophish-powerrat-dcrat/\nauthor: Florian Roth (Nextron Systems)\ndate: 2022-02-21\nmodified: 2025-10-07\ntags:\n- attack.privilege-escalation\n- attack.persistence\n- attack.execution\n- attack.t1053.005\nlogsource:\n  product: windows\n  category: process_creation\ndetection:\n  selection_1_create:\n    Image|endswith: \\schtasks.exe\n    CommandLine|contains|windash: ' /create '\n  selection_1_all_folders:\n    CommandLine|contains:\n    - :\\Perflogs\n    - :\\Users\\All Users\\\n    - :\\Users\\Default\\\n    - :\\Users\\Public\n    - :\\Windows\\Temp\n    - \\AppData\\Local\\\n    - \\AppData\\Roaming\\\n    - '%AppData%'\n    - '%Public%'\n  selection_2_parent:\n    ParentCommandLine|endswith: \\svchost.exe -k netsvcs -p -s Schedule\n  selection_2_some_folders:\n    CommandLine|contains:\n    - :\\Perflogs\n    - :\\Windows\\Temp\n    - \\Users\\Public\n    - '%Public%'\n  filter_optional_other:\n  - ParentCommandLine|contains: unattended.ini\n  - CommandLine|contains: update_task.xml\n  filter_optional_team_viewer:\n    CommandLine|contains: /Create /TN TVInstallRestore /TR\n  filter_optional_avira_install:\n    CommandLine|contains|all:\n    - '/Create /Xml '\n    - \\Temp\\.CR.\n    - \\Avira_Security_Installation.xml\n  filter_optional_avira_other:\n    CommandLine|contains|all:\n    - /Create /F /TN\n    - '/Xml '\n    - \\Temp\\\n    - Avira_\n    CommandLine|contains:\n    - .tmp\\UpdateFallbackTask.xml\n    - .tmp\\WatchdogServiceControlManagerTimeout.xml\n    - .tmp\\SystrayAutostart.xml\n    - .tmp\\MaintenanceTask.xml\n  filter_optional_klite_codec:\n    CommandLine|contains|all:\n    - \\Temp\\\n    - '/Create /TN \"klcp_update\" /XML '\n    - \\klcp_update_task.xml\n  condition: ( all of selection_1_* or all of selection_2_* ) and not 1 of filter_optional_*\nfalsepositives:\n- Benign scheduled tasks creations or executions that happen often during software\n  installations\n- Software that uses the AppData folder and scheduled tasks to update the software\n  in the AppData folders\nlevel: medium\n",
            "vulnerabilities": []
        }
    ]
}