# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_schtasks_one_time_only_midnight_task.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_schtasks_one_time_only_midnight_task.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "pH-T (Nextron Systems)",
                "date": "2022-07-15",
                "title": "Uncommon One Time Only Scheduled Task At 00:00",
                "id": "970823b7-273b-460a-8afc-3a6811998529"
            },
            "rule_text": "title: Uncommon One Time Only Scheduled Task At 00:00\nid: 970823b7-273b-460a-8afc-3a6811998529\nstatus: test\ndescription: Detects scheduled task creation events that include suspicious actions,\n  and is run once at 00:00\nreferences:\n- https://www.trendmicro.com/vinfo/us/security/news/ransomware-spotlight/ransomware-spotlight-blackbyte\nauthor: pH-T (Nextron Systems)\ndate: 2022-07-15\nmodified: 2023-02-03\ntags:\n- attack.execution\n- attack.persistence\n- attack.privilege-escalation\n- attack.t1053.005\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_img:\n  - Image|contains: \\schtasks.exe\n  - OriginalFileName: schtasks.exe\n  selection_cli:\n    CommandLine|contains:\n    - wscript\n    - vbscript\n    - cscript\n    - 'wmic '\n    - wmic.exe\n    - regsvr32.exe\n    - powershell\n    - \\AppData\\\n  selection_time:\n    CommandLine|contains|all:\n    - once\n    - 00:00\n  condition: all of selection_*\nfalsepositives:\n- Software installation\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}