# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_schtasks_folder_combos.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_schtasks_folder_combos.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Florian Roth (Nextron Systems)",
                "date": "2022-04-15",
                "title": "Schtasks From Suspicious Folders",
                "id": "8a8379b8-780b-4dbf-b1e9-31c8d112fefb"
            },
            "rule_text": "title: Schtasks From Suspicious Folders\nid: 8a8379b8-780b-4dbf-b1e9-31c8d112fefb\nstatus: test\ndescription: Detects scheduled task creations that have suspicious action command\n  and folder combinations\nreferences:\n- https://symantec-enterprise-blogs.security.com/blogs/threat-intelligence/lazarus-dream-job-chemical\nauthor: Florian Roth (Nextron Systems)\ndate: 2022-04-15\nmodified: 2022-11-18\ntags:\n- attack.privilege-escalation\n- attack.persistence\n- attack.execution\n- attack.t1053.005\nlogsource:\n  product: windows\n  category: process_creation\ndetection:\n  selection_img:\n  - Image|endswith: \\schtasks.exe\n  - OriginalFileName: schtasks.exe\n  selection_create:\n    CommandLine|contains: ' /create '\n  selection_command:\n    CommandLine|contains:\n    - powershell\n    - pwsh\n    - 'cmd /c '\n    - 'cmd /k '\n    - 'cmd /r '\n    - 'cmd.exe /c '\n    - 'cmd.exe /k '\n    - 'cmd.exe /r '\n  selection_all_folders:\n    CommandLine|contains:\n    - C:\\ProgramData\\\n    - '%ProgramData%'\n  condition: all of selection_*\nfalsepositives:\n- Unknown\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}