# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_sc_service_path_modification.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_sc_service_path_modification.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Victor Sergeev, oscd.community, Nasreddine Bencherchali (Nextron Systems)",
                "date": "2019-10-21",
                "title": "Suspicious Service Path Modification",
                "id": "138d3531-8793-4f50-a2cd-f291b2863d78"
            },
            "rule_text": "title: Suspicious Service Path Modification\nid: 138d3531-8793-4f50-a2cd-f291b2863d78\nstatus: test\ndescription: Detects service path modification via the \"sc\" binary to a suspicious\n  command or path\nreferences:\n- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1543.003/T1543.003.md\n- https://web.archive.org/web/20180331144337/https://www.fireeye.com/blog/threat-research/2018/03/sanny-malware-delivery-method-updated-in-recently-observed-attacks.html\nauthor: Victor Sergeev, oscd.community, Nasreddine Bencherchali (Nextron Systems)\ndate: 2019-10-21\nmodified: 2022-11-18\ntags:\n- attack.persistence\n- attack.privilege-escalation\n- attack.t1543.003\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    Image|endswith: \\sc.exe\n    CommandLine|contains|all:\n    - config\n    - binPath\n    CommandLine|contains:\n    - powershell\n    - 'cmd '\n    - mshta\n    - wscript\n    - cscript\n    - rundll32\n    - svchost\n    - dllhost\n    - cmd.exe /c\n    - cmd.exe /k\n    - cmd.exe /r\n    - cmd /c\n    - cmd /k\n    - cmd /r\n    - C:\\Users\\Public\n    - \\Downloads\\\n    - \\Desktop\\\n    - \\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\\n    - C:\\Windows\\TEMP\\\n    - \\AppData\\Local\\Temp\n  condition: selection\nfalsepositives:\n- Unlikely\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}