# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_reg_direct_asep_registry_keys_modification.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_reg_direct_asep_registry_keys_modification.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Victor Sergeev, Daniil Yugoslavskiy, oscd.community, Swachchhanda Shrawan Poudel (Nextron Systems)",
                "date": "2019-10-25",
                "title": "Direct Autorun Keys Modification",
                "id": "24357373-078f-44ed-9ac4-6d334a668a11"
            },
            "rule_text": "title: Direct Autorun Keys Modification\nid: 24357373-078f-44ed-9ac4-6d334a668a11\nstatus: test\ndescription: Detects direct modification of autostart extensibility point (ASEP) in\n  registry using reg.exe.\nreferences:\n- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md\n- https://github.com/HackTricks-wiki/hacktricks/blob/e4c7b21b8f36c97c35b7c622732b38a189ce18f7/src/windows-hardening/windows-local-privilege-escalation/privilege-escalation-with-autorun-binaries.md\nauthor: Victor Sergeev, Daniil Yugoslavskiy, oscd.community, Swachchhanda Shrawan\n  Poudel (Nextron Systems)\ndate: 2019-10-25\nmodified: 2026-01-05\ntags:\n- attack.privilege-escalation\n- attack.persistence\n- attack.t1547.001\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_img:\n  - Image|endswith: \\reg.exe\n  - OriginalFileName: reg.exe\n  selection_cli_add:\n    CommandLine|contains: add\n  selection_cli_keys:\n    CommandLine|contains:\n    - \\software\\Microsoft\\Windows\\CurrentVersion\\Run\n    - \\software\\WOW6432Node\\Microsoft\\Windows\\CurrentVersion\\Run\n    - \\software\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\n    - \\software\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Userinit\n    - \\software\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell\n    - \\software\\Microsoft\\Windows NT\\CurrentVersion\\Windows\n    - \\system\\CurrentControlSet\\Control\\SafeBoot\\AlternateShell\n  condition: all of selection_*\nfalsepositives:\n- Legitimate software automatically (mostly, during installation) sets up autorun\n  keys for legitimate reasons.\n- Legitimate administrator sets up autorun keys for legitimate reasons.\n- Discord\nlevel: medium\n",
            "vulnerabilities": []
        }
    ]
}