# data/sigma/006d8a4d6e3ea24949907ea9c22f1d5b06467ceaa9e35c49aa44866c854c8901/rules/windows/process_creation/proc_creation_win_schtasks_schedule_type.json
{
    "source_url": "https://github.com/SigmaHQ/sigma/blob/master/rules/windows/process_creation/proc_creation_win_schtasks_schedule_type.yml",
    "rules": [
        {
            "rule_metadata": {
                "status": "test",
                "author": "Nasreddine Bencherchali (Nextron Systems)",
                "date": "2022-09-09",
                "title": "Suspicious Schtasks Schedule Types",
                "id": "24c8392b-aa3c-46b7-a545-43f71657fe98"
            },
            "rule_text": "title: Suspicious Schtasks Schedule Types\nid: 24c8392b-aa3c-46b7-a545-43f71657fe98\nrelated:\n- id: 7a02e22e-b885-4404-b38b-1ddc7e65258a\n  type: similar\nstatus: test\ndescription: Detects scheduled task creations or modification on a suspicious schedule\n  type\nreferences:\n- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-change\n- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create\n- http://blog.talosintelligence.com/2022/09/lazarus-three-rats.html\nauthor: Nasreddine Bencherchali (Nextron Systems)\ndate: 2022-09-09\ntags:\n- attack.privilege-escalation\n- attack.persistence\n- attack.execution\n- attack.t1053.005\nlogsource:\n  product: windows\n  category: process_creation\ndetection:\n  selection_img:\n  - Image|endswith: \\schtasks.exe\n  - OriginalFileName: schtasks.exe\n  selection_time:\n    CommandLine|contains:\n    - ' ONLOGON '\n    - ' ONSTART '\n    - ' ONCE '\n    - ' ONIDLE '\n  filter_privs:\n    CommandLine|contains:\n    - NT AUT\n    - ' SYSTEM'\n    - HIGHEST\n  condition: all of selection_* and not 1 of filter_*\nfalsepositives:\n- Legitimate processes that run at logon. Filter according to your environment\nlevel: high\n",
            "vulnerabilities": []
        }
    ]
}